# ansible/roles/fastapi_app/templates/fastapi-app.service.j2
[Unit]
Description=FastAPI (Uvicorn) service for {{ app_name }}

# Ждём пока сеть поднята полностью и Docker-демон запущен
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
PermissionsStartOnly=true

# 1) Освободить порт, если он вдруг залип после кривого останова
ExecStartPre=/usr/bin/bash -lc '/usr/bin/ss -ltnp | /usr/bin/grep -q ":{{ app_port }} " && /usr/bin/fuser -k {{ app_port }}/tcp || true'

# 2) Дать базе данных в контейнере проснуться.
#    Это критично после ребута, когда uvicorn поднимается быстрее, чем Postgres.
ExecStartPre=/bin/sleep 10

User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_workdir }}

# Прокидываем PYTHONPATH к каталогу приложения
Environment=PYTHONPATH={{ app_workdir }}

# ВАЖНО:
# Ниже через app_env мы прокидываем переменные окружения,
# которые читает FastAPI-приложение через server/config.py.
#
# В частности, сюда должны попасть критичные настройки MQTT-клиента:
#   MQTT_HOST=127.0.0.1
#   MQTT_PORT=1883
#   MQTT_USERNAME=mosquitto-admin
#   MQTT_PASSWORD=***секрет***
#   MQTT_TLS=false
#
# А также сервисные флаги:
#   DEBUG=false                      # в проде off, чтобы скрыть debug-эндпоинты
#   DEVICE_ONLINE_THRESHOLD_S=10     # сколько секунд устройство считается "онлайн" после последнего state
#
# Эти значения задаются в Ansible как словарь app_env и попадают сюда автоматически.
# Ничего не хардкодим в юните, чтобы можно было менять логин/пароль брокера без правки шаблона.
{% for k, v in app_env.items() -%}
Environment={{ k }}={{ v }}
{% endfor %}

# Жёстко целимся в run:app и явно указываем --app-dir
ExecStart={{ app_venv }}/bin/uvicorn --app-dir {{ app_workdir }} run:app --host {{ app_host }} --port {{ app_port }} --workers {{ app_workers }}{% if app_reload %} --reload{% endif %}

ExecStop=/bin/kill -s SIGTERM $MAINPID
KillSignal=SIGINT
TimeoutStopSec=30

# ВАЖНО: хотим автоподъём всегда, а не только "если код возврата плохой"
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
